# startup-scripts/amazon-apiserver.service
[Unit]
Description=AmazonTracker API Server
After=network.target

[Service]
Type=simple
User=node
WorkingDirectory=/opt/amazontracker
EnvironmentFile=-/etc/amazontracker/.env
ExecStart=/usr/bin/node /opt/amazontracker/apiserver.js
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
